# deploy/pebble.service.example
[Unit]
Description=Pebble webmail backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pebble
EnvironmentFile=/opt/pebble/.env
ExecStart=/opt/pebble/target/release/pebble
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
